# systemd/zwhisperd.service
[Unit]
Description=zwhisper recording daemon
After=pipewire.service wireplumber.service
Requires=pipewire.service

[Service]
Type=dbus
BusName=cz.zajca.Zwhisper1
ExecStart=/usr/bin/zwhisperd
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
